# examples/potentiometer_4x.service
[Unit]
Description=Maestra 4x Potentiometer Bridge
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=jordan
WorkingDirectory=/home/jordan/Dev/maestra-core/examples
ExecStart=/home/jordan/Dev/maestra-core/examples/.venv/bin/python examples/potentiometer_4x.py --entity-slug sigint --api-url http://192.168.128.115:8080
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
